# launchd/com.example.seo-onsite-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.seo-onsite-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd "/path/to/the-four-systems" &amp;&amp; ./coordinator.sh onsite-audit</string>
    </array>

    <!-- Weekly Monday 09:17. Lighthouse + on-page audit on the URLs in context/audit-urls.txt. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/seo-onsite-audit.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/seo-onsite-audit.stderr</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Applications/cmux.app/Contents/Resources/bin:/path/to/.local/bin:/path/to/.npm-global/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
